# api/v1/workers/lendro-worker.conf
; Supervisor config for Lendro background workers
;
; Install supervisor: sudo apt install supervisor
; Copy this file to: /etc/supervisor/conf.d/lendro.conf
; Then run: sudo supervisorctl reread && sudo supervisorctl update
;
; Supervisor ensures the workers restart automatically if they crash.

[program:lendro-tx-worker]
command=php /var/www/html/api/v1/workers/process_transactions.php
directory=/var/www/html
autostart=true
autorestart=true
startretries=5
numprocs=2                          ; run 2 worker processes in parallel
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/lendro-worker.log
stderr_logfile=/var/log/lendro-worker-error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
user=www-data                       ; run as the web server user

[program:lendro-recon-worker]
command=php /var/www/html/api/v1/cronjob/reconcile_transactions.php
directory=/var/www/html
autostart=true
autorestart=true
startretries=5
numprocs=1                          ; one reconciler is enough
stdout_logfile=/var/log/lendro-recon.log
stderr_logfile=/var/log/lendro-recon-error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
user=www-data
